# napims-360/napims-360-master/conf/supervisord.napims360.conf
[program:napims360_app]
command=/usr/local/bin/gunicorn napims360.asgi:application -k uvicorn.workers.UvicornWorker --bind :8000
directory=/code/
autostart=true
autorestart=true
stdout_logfile=/var/log/napims360_app/napims360-app.log
redirect_stderr=True
user=napims360
